# etc_systemd_system/cap-part-cleanup.service
[Unit]
After=network.target

[Service]
ExecStart=/srv/scripts/part_cleanup.sh

[Install]
WantedBy=default.target
