# etc_systemd_system/cap2ascii.service
[Unit]
After=network.target

[Service]
ExecStart=/srv/scripts/capPay2ascii.sh

[Install]
WantedBy=default.target
