# etc_systemd_system/pcap-start.service
[Unit]
After=network.target

[Service]
ExecStart=/srv/scripts/startSniff.sh

[Install]
WantedBy=default.target
